# ansible/roles/go/files/profile.d/golang.sh
export GOPATH="$HOME/go"
export GOBIN="$GOPATH/bin"
mkdir -p $GOBIN
export PATH="$PATH:/usr/local/go/bin:$GOPATH:$GOBIN"
